# rdc/docker-cron/crontabs
# REDCap Cron Job (runs every minute)
* * * * *   .{REDCAP_WEBROOT_PATH}cron.php --spider >/dev/null 2>&1
0 * * * *   /usr/sbin/logrotate /etc/logrotate.conf