# rdc/docker-cron/crontabs
* *   * * *   wget webREDCAP_WEBROOT_PATHcron.php --spider >/dev/null 2>&1
#* */3 * * *   /usr/sbin/logrotate /etc/logrotate.conf

# REDCap Cron Job (runs every minute)
* * * * * /usr/local/bin/php /var/www/html/cron.php > /dev/null
